# vbk-supervisord.conf
; VietByte AI — supervisord program config
; Copy file này vào /etc/supervisor/conf.d/vbk-ai.conf rồi: supervisorctl reread && supervisorctl update
;
; 4 service độc lập, KHÔNG có gateway proxy:
;   vbk-vllm-translator (port 9001) — vLLM Qwen3-14B-AWQ, internal
;   vbk-translate       (port 9002) — /v1/translate, /v1/json, /v1/dict (gọi vLLM)
;   vbk-ocr             (port 9003) — /v1/ocr (PaddleOCR PP-OCRv5 GPU)
;   vbk-tts             (port 9004) — /v1/tts (Chatterbox Multilingual 0.5B GPU)

[program:vbk-vllm-translator]
command=/workspace/vbk-ai-server/scripts/run_vllm_translator.sh
directory=/workspace/vbk-ai-server
user=root
autostart=true
autorestart=true
startsecs=30
startretries=3
stopwaitsecs=20
stopasgroup=true
killasgroup=true
stdout_logfile=/workspace/vbk-ai-server/logs/vllm-translator.log
stdout_logfile_maxbytes=20MB
stdout_logfile_backups=3
stderr_logfile=/workspace/vbk-ai-server/logs/vllm-translator-err.log
stderr_logfile_maxbytes=20MB
stderr_logfile_backups=3
priority=10

[program:vbk-ocr]
command=/workspace/vbk-ai-server/scripts/run_ocr_service.sh
directory=/workspace/vbk-ai-server/ocr_service
user=root
autostart=true
autorestart=true
startsecs=15
startretries=3
stopwaitsecs=15
stopasgroup=true
killasgroup=true
stdout_logfile=/workspace/vbk-ai-server/logs/ocr.log
stdout_logfile_maxbytes=20MB
stdout_logfile_backups=3
stderr_logfile=/workspace/vbk-ai-server/logs/ocr-err.log
stderr_logfile_maxbytes=20MB
stderr_logfile_backups=3
priority=15

[program:vbk-tts]
command=/workspace/vbk-ai-server/scripts/run_tts_service.sh
directory=/workspace/vbk-ai-server/tts_service
user=root
autostart=true
autorestart=true
startsecs=60          ; Chatterbox warm-up lần đầu 60-120s (load model 2GB + CUDA kernel compile)
startretries=3
stopwaitsecs=20
stopasgroup=true
killasgroup=true
stdout_logfile=/workspace/vbk-ai-server/logs/tts.log
stdout_logfile_maxbytes=20MB
stdout_logfile_backups=3
stderr_logfile=/workspace/vbk-ai-server/logs/tts-err.log
stderr_logfile_maxbytes=20MB
stderr_logfile_backups=3
priority=18

[program:vbk-translate]
command=/workspace/vbk-ai-server/scripts/run_translate_service.sh
directory=/workspace/vbk-ai-server/translate_service
user=root
autostart=true
autorestart=true
startsecs=5
startretries=3
stopwaitsecs=10
stopasgroup=true
killasgroup=true
stdout_logfile=/workspace/vbk-ai-server/logs/translate.log
stdout_logfile_maxbytes=20MB
stdout_logfile_backups=3
stderr_logfile=/workspace/vbk-ai-server/logs/translate-err.log
stderr_logfile_maxbytes=20MB
stderr_logfile_backups=3
priority=20

[group:vbk-ai]
programs=vbk-vllm-translator,vbk-ocr,vbk-tts,vbk-translate
