# vbk-supervisord.conf
; VietByte AI — supervisord program config
; Copy file này vào /etc/supervisor/conf.d/vbk-ai.conf rồi: supervisorctl reread && supervisorctl update
;
; Stack: vLLM translator (Qwen3-14B-AWQ) + OCR service (PaddleOCR) + API gateway.
; OCR tách riêng → API gateway không bị block bởi OCR-heavy.

[program:vbk-vllm-translator]
command=/workspace/vbk-ai-server/scripts/run_vllm_translator.sh
directory=/workspace/vbk-ai-server
user=root
autostart=true
autorestart=true
startsecs=30
startretries=3
stopwaitsecs=20
stopasgroup=true
killasgroup=true
stdout_logfile=/workspace/vbk-ai-server/logs/vllm-translator.log
stdout_logfile_maxbytes=20MB
stdout_logfile_backups=3
stderr_logfile=/workspace/vbk-ai-server/logs/vllm-translator-err.log
stderr_logfile_maxbytes=20MB
stderr_logfile_backups=3
priority=10

[program:vbk-ocr]
command=/workspace/vbk-ai-server/scripts/run_ocr_service.sh
directory=/workspace/vbk-ai-server/ocr_service
user=root
autostart=true
autorestart=true
startsecs=15
startretries=3
stopwaitsecs=15
stopasgroup=true
killasgroup=true
stdout_logfile=/workspace/vbk-ai-server/logs/ocr.log
stdout_logfile_maxbytes=20MB
stdout_logfile_backups=3
stderr_logfile=/workspace/vbk-ai-server/logs/ocr-err.log
stderr_logfile_maxbytes=20MB
stderr_logfile_backups=3
priority=15

[program:vbk-api]
command=/workspace/vbk-ai-server/scripts/run_api.sh
directory=/workspace/vbk-ai-server/api
user=root
autostart=true
autorestart=true
startsecs=5
startretries=3
stopwaitsecs=10
stopasgroup=true
killasgroup=true
stdout_logfile=/workspace/vbk-ai-server/logs/api.log
stdout_logfile_maxbytes=20MB
stdout_logfile_backups=3
stderr_logfile=/workspace/vbk-ai-server/logs/api-err.log
stderr_logfile_maxbytes=20MB
stderr_logfile_backups=3
priority=20

[group:vbk-ai]
programs=vbk-vllm-translator,vbk-ocr,vbk-api
